# deploy/systemd/futuresml-worker.service
[Unit]
Description=FuturesML ML experiment worker
After=network.target postgresql.service

[Service]
WorkingDirectory=/root/FuturesML
EnvironmentFile=/root/FuturesML/.env
ExecStart=/root/FuturesML/.venv/bin/python -m app.workers.ml_experiment_worker
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
